# deploy/brewsignal.service
[Unit]
Description=BrewSignal - Beer Fermentation Dashboard
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/brewsignal
Environment="PATH=/opt/brewsignal/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/brewsignal/.venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5

# Security hardening (NoNewPrivileges disabled to allow sudo for system controls)
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/brewsignal/data
PrivateTmp=true

# Bluetooth access
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=brewsignal

[Install]
WantedBy=multi-user.target
